# profile.d/default
backup_modules="
  filesystem
"

backup_environ="
  RESTIC_PASSWORD=changeme
  RESTIC_REPOSITORY=/mnt/backups
"

backup_basedir="/var/backups"
